# services/imputation-server/packer/node_exporter.conf
description "Start node_exporter"
start on filesystem or runlevel [2345]
stop on shutdown
respawn
exec /usr/local/bin/node_exporter
